# mikrom-network/debian/lib/systemd/system/mikrom-network.service
[Unit]
Description=Mikrom Network Mesh Orchestrator
After=network.target nats-server.service
Wants=network.target

[Service]
Type=simple
User=root
EnvironmentFile=/etc/mikrom/network.env
ExecStart=/usr/bin/mikrom-network
Restart=always
RestartSec=5
CapabilityBoundingSet=CAP_NET_ADMIN
AmbientCapabilities=CAP_NET_ADMIN
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
