# mikrom-agent/debian/mikrom-agent.service
[Unit]
Description=Mikrom Agent
After=network.target

[Service]
Type=simple
EnvironmentFile=-/etc/mikrom/agent.env
ExecStart=/usr/bin/mikrom-agent
Restart=always
RestartSec=5
User=root
# El agente necesita acceso a /dev/kvm para Firecracker
Group=root
# Deja que los hijos sobrevivan al reinicio del servicio; el agente gestiona
# Firecracker como procesos aparte.
KillMode=process

[Install]
WantedBy=multi-user.target
